# .config/systemd/user/getNumArchUpd.service
[Unit]
Description=Service to get the number of Pacman updates available

[Service]
ExecStart=/home/papa/bin/xmbGetNumArchUpd

