# .config/systemd/user/getNumArchUpd.timer
[Unit]
Description=Check number of Arch pacman updates daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

